# infra/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ── ambot Shared Engine (T1 + T2 clients) ────────────────────────────────────
[program:ambot_shared]
command=python -m ambot.main
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
stdout_logfile=/app/logs/ambot_shared.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/ambot_shared_err.log
environment=
    ENV="production",
    LOG_LEVEL="INFO"

# ── Web API ───────────────────────────────────────────────────────────────────
[program:web]
command=uvicorn web.main:app --host 0.0.0.0 --port 8000 --workers 2
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/app/logs/web.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/web_err.log
